# systemd/mas004-esp32-plc-bridge.service
[Unit]
Description=MAS-004 ESP32 PLC Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/MAS-004_ESP32-PLC-Bridge
ExecStart=/opt/MAS-004_ESP32-PLC-Bridge/.venv/bin/python -m mas004_esp32_plc_bridge --config /etc/mas004_esp32_plc_bridge/config.json
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
